# bin/webmonitor.service
[Unit]
Description=Webmonitor - daily text diffs of websites

[Service]
ExecStart=/home/wdm/bin/webmonitor.sh

[Install]
WantedBy=multi-user.target
