# rc.local
#!/bin/sh -e
#
# rc.local
#
# This script is executed at the end of each multiuser runlevel.
# Make sure that the script will "exit 0" on success or any other
# value on error.

# Printing.
/usr/sbin/cupsd

# USB Devices
# See https://github.com/dnschneid/crouton/wiki/udev:-manage-inserted-devices
mkdir -p /run/udev/rules.d
cp /etc/udev/rules.d/*.rules /run/udev/rules.d
udevadm control --reload

exit 0
